# deploy_tools/supervisord.conf.j2
[program:gunicorn]
directory=/home/oliver/sites/{{ host }}/source
command=/home/oliver/sites/{{ host }}/source/gunicorn-supervisor.conf
user=oliver
umask=022
