# deploy_tools/supervisord.template.conf
[program:gunicorn]
directory=/home/oliver/sites/SITENAME/source
command=/home/oliver/sites/SITENAME/source/gunicorn-supervisor.conf
user=oliver
umask=022
